# app-portage/creideiki-gentoo-utils/files/update-portage.cron
#!/bin/sh

( emaint sync -a && eix-update ) > /dev/null 2> /dev/null
glsa-check --quiet --nocolor --list affected
eclean --nocolor --quiet --time-limit=1m packages
eclean --nocolor --quiet --deep distfiles
fstrim --all --verbose
